# Homework_3/appboot.service
[Unit]
Description=Start app

[Service]
Type=oneshot
ExecStart=/home/vagrant/script

[Install]
WantedBy=multi-user.target